# luci-app-adguardhome/root/etc/uci-defaults/40_luci-AdGuardHome
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@AdGuardHome[-1]
	add ucitrack AdGuardHome
	set ucitrack.@AdGuardHome[-1].init=AdGuardHome
	commit ucitrack
EOF

# 注册 Firewall Include
# 检查是否已存在
if ! uci -q get firewall.adguardhome_include >/dev/null; then
	uci -q batch <<-EOF
		set firewall.adguardhome_include=include
		set firewall.adguardhome_include.type='script'
		set firewall.adguardhome_include.path='/usr/share/AdGuardHome/firewall.include'
		set firewall.adguardhome_include.family='any'
		set firewall.adguardhome_include.reload='1'
		commit firewall
	EOF
fi

# 清理旧版本脚本可能残留的动态防火墙规则
# 旧版脚本使用 name='AdGuard Home DNS'
old_rules=$(uci show firewall | grep "name='AdGuard Home DNS'" | cut -d. -f2)
for rule in $old_rules; do
	uci -q delete firewall."$rule"
done
[ -n "$old_rules" ] && uci commit firewall

# 清理旧的计划任务 (tailto.sh 等)
sed -i '/tailto.sh/d' /etc/crontabs/root 2>/dev/null
sed -i '/waitnet.sh/d' /etc/crontabs/root 2>/dev/null
sed -i '/watchconfig.sh/d' /etc/crontabs/root 2>/dev/null

exit 0