# luci-app-adguardhome/root/etc/init.d/AdGuardHome
#!/bin/sh /etc/rc.common

USE_PROCD=1

START=95
STOP=01

CONFIGURATION=AdGuardHome
CRON_FILE=/etc/crontabs/root
GFWSET="gfwlist"
EXTRA_COMMANDS="do_redirect testbackup test_crontab force_reload isrunning"
EXTRA_HELP="        do_redirect 0 or 1\
        testbackup backup or restore\
		test_crontab
		force_reload
		isrunning"
start_forward_dnsmasq()
{
	addr="127.0.0.1#$AdGuardHome_PORT"
	OLD_SERVER="`uci get dhcp.@dnsmasq[0].server 2>/dev/null`"
	echo $OLD_SERVER | grep "^$addr" >/dev/null 2>&1
	if [ $? -eq 0 ]; then
		return
	fi
	uci delete dhcp.@dnsmasq[0].server 2>/dev/null
	uci add_list dhcp.@dnsmasq[0].server=$addr
	for server in $OLD_SERVER; do
		if [ "$server" = "$addr" ]; then
			continue
		fi
		uci add_list dhcp.@dnsmasq[0].server=$server
	done
	uci delete dhcp.@dnsmasq[0].resolvfile 2>/dev/null
	uci set dhcp.@dnsmasq[0].noresolv=1
	uci commit dhcp
	/etc/init.d/dnsmasq restart
}

stop_forward_dnsmasq()
{
	local OLD_PORT="$1"
	addr="127.0.0.1#$OLD_PORT"
	OLD_SERVER="`uci get dhcp.@dnsmasq[0].server 2>/dev/null`"
	echo $OLD_SERVER | grep "^$addr" >/dev/null 2>&1
	if [ $? -ne 0 ]; then
		return
	fi

	uci del_list dhcp.@dnsmasq[0].server=$addr 2>/dev/null
	rcauto="/tmp/resolv.conf.auto"
	if [ -e "/tmp/resolv.conf.d/resolv.conf.auto" ]; then
		rcauto="/tmp/resolv.conf.d/resolv.conf.auto"
	fi
	addrlist="`uci get dhcp.@dnsmasq[0].server 2>/dev/null`"
	if [ -z "$addrlist" ] ; then
		uci set dhcp.@dnsmasq[0].resolvfile="$rcauto" 2>/dev/null
		uci delete dhcp.@dnsmasq[0].noresolv 2>/dev/null
	fi
	uci commit dhcp
	/etc/init.d/dnsmasq restart
}

set_firewall()
{
	uci set firewall.adguardhome_redirect="redirect"
	uci set firewall.adguardhome_redirect.target="DNAT"
	uci set firewall.adguardhome_redirect.name="AdGuard Home"
	uci set firewall.adguardhome_redirect.src="lan"
	uci set firewall.adguardhome_redirect.family="any"
	uci set firewall.adguardhome_redirect.src_dport="53"
	uci set firewall.adguardhome_redirect.dest_port="$AdGuardHome_PORT"
	uci commit firewall
	/etc/init.d/firewall reload
	echo "Firewall rules updated."
}

clear_firewall()
{
	redirects=$(uci show firewall | grep "firewall.adguardhome_redirect")
	if [ -z "$redirects" ]; then
		echo "No redirect rules found."
	else
		uci delete firewall.adguardhome_redirect
		uci commit firewall
		/etc/init.d/firewall reload
		echo "Deleted redirect rule: firewall.adguardhome_redirect"
	fi
}

service_triggers() {
	procd_add_reload_trigger "$CONFIGURATION"
	[ "$(uci get AdGuardHome.AdGuardHome.redirect)" == "redirect" ] && procd_add_reload_trigger firewall
}
isrunning(){
	config_load "${CONFIGURATION}"
	_isrunning
	local r=$?
	([ "$r" == "0" ] && echo "running") || ([ "$r" == "1" ] && echo "not run" ) || echo "no bin"
	return $r
}
_isrunning(){
	config_get binpath $CONFIGURATION binpath "/usr/bin/AdGuardHome/AdGuardHome"
	[ ! -f "$binpath" ] && return 2
	pgrep $binpath 2>&1 >/dev/null && return 0
	return 1
}
force_reload(){
	config_load "${CONFIGURATION}"
	_isrunning && procd_send_signal "$CONFIGURATION" || start
}

rm_port53()
{
	local AdGuardHome_PORT=$(config_editor "dns.port" "" "$configpath" "1")
	dnsmasq_port=$(uci get dhcp.@dnsmasq[0].port 2>/dev/null)
	if [ -z "$dnsmasq_port" ]; then
		dnsmasq_port="53"
	fi
	if [ "$dnsmasq_port" == "$AdGuardHome_PORT" ]; then
		if [ "$dnsmasq_port" == "53" ]; then
			dnsmasq_port="1745"
		fi
	elif [ "$dnsmasq_port" == "53" ]; then
		return
	fi
	config_editor "dns.port" "$dnsmasq_port" "$configpath"
	uci set dhcp.@dnsmasq[0].port="53"
	uci commit dhcp
	/etc/init.d/dnsmasq restart
}

use_port53()
{
	local AdGuardHome_PORT=$(config_editor "dns.port" "" "$configpath" "1")
	dnsmasq_port=$(uci get dhcp.@dnsmasq[0].port 2>/dev/null)
	if [ -z "$dnsmasq_port" ]; then
		dnsmasq_port="53"
	fi
	if [ "$dnsmasq_port" == "$AdGuardHome_PORT" ]; then
		if [ "$dnsmasq_port" == "53" ]; then
			AdGuardHome_PORT="1745"
		fi
	elif [ "$AdGuardHome_PORT" == "53" ]; then
		return
	fi
	config_editor "dns.port" "53" "$configpath"
	uci set dhcp.@dnsmasq[0].port="$AdGuardHome_PORT"
	uci commit dhcp
	/etc/init.d/dnsmasq restart
}
do_redirect()
{
	config_load "${CONFIGURATION}"
	_do_redirect $1
}
_do_redirect()
{
	local section="$CONFIGURATION"
	enabled=$1
	config_get configpath $CONFIGURATION configpath "/etc/AdGuardHome.yaml"
	AdGuardHome_PORT=$(config_editor "dns.port" "" "$configpath" "1")
	if [ -z "$AdGuardHome_PORT" ]; then
		AdGuardHome_PORT="0"
	fi
	config_get "redirect" "$section" "redirect" "none"
	config_get "old_redirect" "$section" "old_redirect" "none"
	config_get "old_port" "$section" "old_port" "0"
	config_get "old_enabled" "$section" "old_enabled" "0"
	uci get dhcp.@dnsmasq[0].port >/dev/null 2>&1 || uci set dhcp.@dnsmasq[0].port="53" >/dev/null 2>&1
	if [ "$old_enabled" = "1" -a "$old_redirect" == "exchange" ]; then
		AdGuardHome_PORT=$(uci get dhcp.@dnsmasq[0].port 2>/dev/null)
	fi

	if [ "$old_redirect" != "$redirect" ] || [ "$old_port" != "$AdGuardHome_PORT" ] || [ "$old_enabled" = "1" -a "$enabled" = "0" ]; then
		if [ "$old_redirect" != "none" ]; then
			if [  "$old_redirect" == "redirect" -a "$old_port" != "0" ]; then
				clear_firewall
			elif [ "$old_redirect" == "dnsmasq-upstream" ]; then
				stop_forward_dnsmasq "$old_port"
			elif [ "$old_redirect" == "exchange" ]; then
				rm_port53
			fi
		fi
	elif [ "$old_enabled" = "1" -a "$enabled" = "1" ]; then
		if [  "$old_redirect" == "redirect" -a "$old_port" != "0" ]; then
			clear_firewall
		fi
	fi
	uci delete AdGuardHome.@AdGuardHome[0].old_redirect 2>/dev/null
	uci delete AdGuardHome.@AdGuardHome[0].old_port 2>/dev/null
	uci delete AdGuardHome.@AdGuardHome[0].old_enabled 2>/dev/null
	uci add_list AdGuardHome.@AdGuardHome[0].old_redirect="$redirect" 2>/dev/null
	uci add_list AdGuardHome.@AdGuardHome[0].old_port="$AdGuardHome_PORT" 2>/dev/null
	uci add_list AdGuardHome.@AdGuardHome[0].old_enabled="$enabled" 2>/dev/null
	uci commit AdGuardHome

	if [ "$enabled" == "0" ] || [ "$AdGuardHome_PORT" == "0" ]; then
		echo -n "0">/var/run/AdG_redir
		return 1
	fi
	redirected=1
	if [ "$redirect" = "redirect" ]; then
		set_firewall
	elif [ "$redirect" = "dnsmasq-upstream" ]; then
		start_forward_dnsmasq
	elif [ "$redirect" == "exchange" -a "$(uci get dhcp.@dnsmasq[0].port 2>/dev/null)" == "53" ]; then
		use_port53
	else
		redirected=0
	fi

	echo -n "$redirected">/var/run/AdG_redir
}
get_filesystem()
{
# print out path filesystem
    echo $1 | awk '
    BEGIN{
    while (("mount"| getline ret) > 0)
    {
    split(ret,d);
    fs[d[3]]=d[5];
    m=index(d[1],":")
    if (m==0)
    {
        pt[d[3]]=d[1]
    }else{
        pt[d[3]]=substr(d[1],m+1)
    }}}{
    split($0,d,"/");
    if ("/" in fs)
    {
    result1=fs["/"];
    }
    if ("/" in pt)
    {
    result2=pt["/"];
    }
    for (i=2;i<=length(d);i++)
    {
       p[i]=p[i-1]"/"d[i];
       if (p[i] in fs)
        {
        result1=fs[p[i]];
        result2=pt[p[i]];
        }
    }
    if (result2 in fs){
        result=fs[result2]}
    else{
        result=result1}
    print(result);}'
}

config_editor()
{
	[ -n "$3" -a -f $3 ] || return
	awk -v yaml="$1" -v value="$2" -v file="$3" -v ro="$4" '
	BEGIN{split(yaml,part,"\.");s="";i=1;l=length(part);}
    {
        if (match($0,s""part[i]":"))
        {
            if (i==l)
            {
                split($0,t,": ");
				if (ro==""){
				system("sed -i '\''"FNR"c \\"t[1]": "value"'\'' "file);
				}else{
				print(t[2]);
				}
				exit;
            }
            s=s"[- ]{2}";
            i++;
        }
    }' $3
}

boot_service() {
	config_load "${CONFIGURATION}"
	config_get waitonboot $CONFIGURATION waitonboot "0"
	config_get_bool enabled $CONFIGURATION enabled 0
	config_get binpath $CONFIGURATION binpath "/usr/bin/AdGuardHome/AdGuardHome"
	[ -f "$binpath" ] && start_service
	if [ "$enabled" == "1" ] && [ "$waitonboot" == "1" ]; then
		procd_open_instance "waitnet"
		procd_set_param command "/usr/share/AdGuardHome/waitnet.sh"
		procd_close_instance
		echo "No net start pinging"
	fi
}

testbackup(){
	config_load "${CONFIGURATION}"
	if [ "$1" == "backup" ]; then
		backup
	elif [ "$1" == "restore" ]; then
		restore
	fi
}
restore()
{
	config_get workdir $CONFIGURATION workdir "/usr/bin/AdGuardHome"
	config_get backupwdpath $CONFIGURATION backupwdpath "/usr/bin/AdGuardHome"
	cp -u -r -f $backupwdpath/data $workdir
}
backup() {
	config_get backupwdpath $CONFIGURATION backupwdpath "/usr/bin/AdGuardHome"
	mkdir -p $backupwdpath/data
	config_get workdir $CONFIGURATION workdir "/usr/bin/AdGuardHome"
	config_get backupfile $CONFIGURATION backupfile ""
	for one in $backupfile;
	do
		while :
		do
			if [ -d "$backupwdpath/data/$one" ]; then
				cpret=$(cp -u -r -f $workdir/data/$one $backupwdpath/data 2>&1)
			else
				cpret=$(cp -u -r -f $workdir/data/$one $backupwdpath/data/$one 2>&1)
			fi
			echo "$cpret"
			echo "$cpret" | grep "no space left on device"
			if [ "$?" == "0" ]; then
				echo "Disk full, deleting log and retrying"
				del_querylog && continue
				rm -f -r $backupwdpath/data/filters
				rm -f -r $workdir/data/filters && continue
				echo "Backup failed"
			fi
			break
		done
	done
}

_expand_upprotect_path() {
    local in="$1"
    local out="$in"

    [ -n "$binpath" ]     && out="${out//\$binpath/$binpath}"
    [ -n "$configpath" ]  && out="${out//\$configpath/$configpath}"
    [ -n "$logfile" ]     && out="${out//\$logfile/$logfile}"
    [ -n "$workdir" ]     && out="${out//\$workdir/$workdir}"

    case "$out" in
        *'$'*)
            logger "AdGuardHome" "upprotect reject unsafe path: $in"
            return 1
        ;;
    esac

    echo "$out"
}

_add_upprotect_file()
{
	local raw="$1"
	local expanded
	# Expand variables in the file path
	expanded="$(_expand_upprotect_path "$raw")" || return 0
	echo "$expanded" >> /lib/upgrade/keep.d/luci-app-adguardhome
}

start_service() {
	# Reading config
	config_load "${CONFIGURATION}"
	# update password
	config_get username $CONFIGURATION username ""
	config_get hashpass $CONFIGURATION hashpass ""
	config_get configpath $CONFIGURATION configpath "/etc/AdGuardHome.yaml"
	if [ -n "$username" ]; then
		config_editor "users.name" "$username" "$configpath"
		uci set $CONFIGURATION.$CONFIGURATION.username=""
	fi
	if [ -n "$hashpass" ]; then
		config_editor "users.password" "$hashpass" "$configpath"
		uci set $CONFIGURATION.$CONFIGURATION.hashpass=""
	fi
	local ipst=0
	config_get_bool enabled $CONFIGURATION enabled 0
	# update crontab
	update_crontab
	echo "luci enable switch=$enabled"
	if [ "$enabled" == "0" ]; then
		_do_redirect 0
		return
	fi
	#what need to do before reload
	config_get workdir $CONFIGURATION workdir "/usr/bin/AdGuardHome"
	grep -q "ipset.txt" $configpath 2>/dev/null && ipst=1
	if [ $ipst -eq 1 ];then
		ipset list $GFWSET >/dev/null 2>&1 || ipset create $GFWSET hash:ip 2>/dev/null
	fi
	config_get backupfile $CONFIGURATION backupfile ""
	mkdir -p $workdir/data
	if [ -n "$backupfile" ] && [ ! -d "$workdir/data" ]; then
		restore
	fi

	# for overlay data-stk-oo not suppport
	local cwdfs=$(get_filesystem $workdir)
	echo "Workdir $workdir is a $cwdfs filesystem"
	if [ "$cwdfs" == "jffs2" ]; then
		echo "Linking db files to /tmp to avoid FS error"
		logger "AdGuardHome" "warning db redirect to tmp"
		touch $workdir/data/stats.db
		if [ ! -L $workdir/data/stats.db ]; then
		mv -f $workdir/data/stats.db /tmp/stats.db 2>/dev/null
		ln -s /tmp/stats.db $workdir/data/stats.db 2>/dev/null
		fi
		touch $workdir/data/sessions.db
		if [ ! -L $workdir/data/sessions.db ]; then
		mv -f $workdir/data/sessions.db /tmp/sessions.db 2>/dev/null
		ln -s /tmp/sessions.db $workdir/data/sessions.db 2>/dev/null
		fi
	fi

	# hack to save config file when upgrade system
	config_get upprotect $CONFIGURATION upprotect ""
	if [ -n "$upprotect" ]; then
		# Handle both old MultiValue format (space-separated) and new DynamicList format
		echo -n > /lib/upgrade/keep.d/luci-app-adguardhome
		if echo "$upprotect" | grep -q " "; then
			# Old format: space-separated values
			for one in $upprotect; do
				_add_upprotect_file "$one"
			done
		else
			# New format: UCI list, get all values
			config_list_foreach $CONFIGURATION upprotect _add_upprotect_file
		fi
	fi

	config_get binpath $CONFIGURATION binpath "/usr/bin/AdGuardHome/AdGuardHome"
	mkdir -p ${binpath%/*}
	if [ ! -f "$binpath" ]; then
		_do_redirect 0
		/usr/share/AdGuardHome/update_core.sh 2>&1 >/tmp/AdGuardHome_update.log &
		exit 0
	fi

	config_get httpport $CONFIGURATION httpport 3000
	local ADDITIONAL_ARGS=" -c $configpath -w $workdir -p $httpport --no-check-update"

	config_get logfile $CONFIGURATION logfile ""
	if [ -n "$logfile" ]; then
		ADDITIONAL_ARGS="$ADDITIONAL_ARGS -l $logfile"
	fi

	config_get_bool verbose $CONFIGURATION verbose 0
	if [ "$verbose" -eq 1 ]; then
		ADDITIONAL_ARGS="$ADDITIONAL_ARGS -v"
	fi

	procd_open_instance
	SET_TZ="`uci get system.@system[0].timezone 2>/dev/null`"
	if [ -n "$SET_TZ" ]; then
		procd_set_param env TZ="$SET_TZ"
	fi
	procd_set_param respawn ${respawn_threshold:-3600} ${respawn_timeout:-5} ${respawn_retry:-5}
	procd_set_param limits core="unlimited" nofile="65535 65535"
	procd_set_param stderr 1
	procd_set_param command $binpath $ADDITIONAL_ARGS
	procd_set_param file "$configpath" "/etc/hosts" "/etc/config/AdGuardHome"
	procd_close_instance
	if [ -f "$configpath" ]; then
		_do_redirect 1
	else
		_do_redirect 0
		config_get "redirect" $CONFIGURATION "redirect" "none"
		if [ "$redirect" != "none" ]; then
			procd_open_instance "waitconfig"
			procd_set_param command "/usr/share/AdGuardHome/watchconfig.sh"
			procd_close_instance
			echo "No config start watching"
		fi
	fi
	echo "AdGuardHome service enabled"
	(sleep 10 && [ -z "$(pgrep $binpath)" ] && logger "AdGuardHome" "no process in 10s cancel redirect" && _do_redirect 0 )&
}
reload_service()
{
	echo "AdGuardHome reloading"
	start
}

del_querylog(){
	local btarget=$(ls $backupwdpath/data | grep -F "querylog.json" | sort -r | head -n 1)
	local wtarget=$(ls $workdir/data | grep -F "querylog.json" | sort -r | head -n 1)
	if [ "$btarget"x == "$wtarget"x ]; then
		[ -z "$btarget" ] && return 1
		rm -f $workdir/data/$wtarget
		rm -f $backupwdpath/data/$btarget
		return 0
	fi
	if [ "$btarget" \> "$wtarget" ]; then
		rm -f $backupwdpath/data/$btarget
		return 0
	else
		rm -f $workdir/data/$wtarget
		return 0
	fi
}
stop_service()
{
	config_load "${CONFIGURATION}"
	_do_redirect 0
	update_crontab
	if [ "$1" != "nobackup" ]; then
		config_get backupfile $CONFIGURATION backupfile "0"
		if [ -n "$backupfile" ]; then
			backup
		fi
	fi
	echo "AdGuardHome service disabled"
}
boot() {
	rc_procd boot_service "$@"
	if type service_started >/dev/null 2>&1; then
		service_started
	fi
}
test_crontab(){
	config_load "${CONFIGURATION}"
	update_crontab
}
update_crontab(){
	config_get_bool enabled $CONFIGURATION enabled 0
	config_get crontab $CONFIGURATION crontab ""
	local findstr default cronenable replace
	local cronreload=0

	findstr="/usr/share/AdGuardHome/addhost.sh"
	default="0 * * * * /usr/share/AdGuardHome/addhost.sh"
	[ "$enabled" == "0" ] || [ "${crontab//autohost/}" == "$crontab" ] && cronenable=0 || cronenable=1
	crontab_editor
	[ "$cronenable" == "0" ] && /usr/share/AdGuardHome/addhost.sh "del" "noreload" || /usr/share/AdGuardHome/addhost.sh "" "noreload"

	findstr="/usr/share/AdGuardHome/gfw2adg.sh"
	default="30 3 * * * /usr/share/AdGuardHome/gfw2adg.sh"
	[ "$enabled" == "0" ] || [ "${crontab//autogfw/}" == "$crontab" ] && cronenable=0 || cronenable=1
	crontab_editor

	findstr="/usr/share/AdGuardHome/gfwipset2adg.sh"
	default="31 3 * * * /usr/share/AdGuardHome/gfwipset2adg.sh"
	[ "$enabled" == "0" ] || [ "${crontab//autogfwipset/}" == "$crontab" ] && cronenable=0 || cronenable=1
	crontab_editor
	[ "$cronreload" -gt 0 ] && /etc/init.d/cron restart
}
crontab_editor(){
	[ ! -f "$CRON_FILE" ] && mkdir -p ${CRON_FILE%/*} && touch "$CRON_FILE" && chmod 600 "$CRON_FILE"
	#usage input:
	#findstr=
	#default=
	#cronenable=
	#replace="${last//\//\\/}/${now//\//\\/}"
	#output:cronreload:if >1 please /etc/init.d/cron restart manual
	local testline reload
	local line="$(grep "$findstr" $CRON_FILE)"
	[ -n "$replace" ] && [ -n "$line" ] && eval testline="\${line//$replace}" && [ "$testline" != "$line" ] && line="$testline" && reload="1" && replace=""
	if [ "${line:0:1}" != "#" ]; then
		if [ $cronenable -eq 1 ]; then
			[ -z "$line" ] && line="$default" && reload="1"
			if [ -n "$reload" ]; then
				sed -i "\,$findstr,d" $CRON_FILE
				echo "$line" >> $CRON_FILE
				cronreload=$((cronreload+1))
			fi
		elif [ -n "$line" ]; then
			sed -i "\,$findstr,d" $CRON_FILE
			echo "#$line" >> $CRON_FILE
			cronreload=$((cronreload+1))
		fi
	else
		if [ $cronenable -eq 1 ]; then
			sed -i "\,$findstr,d" $CRON_FILE
			echo "${line:1}" >> $CRON_FILE
			cronreload=$((cronreload+1))
		elif [ -z "$reload" ]; then
			sed -i "\,$findstr,d" $CRON_FILE
			echo "$line" >> $CRON_FILE
		fi
	fi
}
